# scripts/deploy/brain-processing-audit.timer
[Unit]
Description=Run AI Brain Processing integrity audit every six hours

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
Persistent=true
RandomizedDelaySec=5m
Unit=brain-processing-audit.service

[Install]
WantedBy=timers.target
