# scripts/deploy/brain-notebooklm-retention.timer
[Unit]
Description=Run the AI Brain NotebookLM durable retention fallback every minute

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m
AccuracySec=1s
RandomizedDelaySec=0
Persistent=true
Unit=brain-notebooklm-retention.service

[Install]
WantedBy=timers.target
