# scripts/deploy/brain-recall-sync.service
[Unit]
Description=AI Brain Recall daily snapshot import
After=network-online.target brain.service
Wants=network-online.target

[Service]
Type=oneshot
User=brain-recall
Group=brain-data
WorkingDirectory=/opt/brain
EnvironmentFile=/etc/brain/.env
Environment=BRAIN_RECALL_TRIGGER=automatic
Environment=BRAIN_RECALL_MANUAL_SYNC_MODE=0
Environment=BRAIN_RECALL_BACKUP_STAGING_DIR=/run/brain-recall-backup-staging
Environment=BRAIN_SCRUB_RUNTIME_ROOT=/opt/brain/current
LoadCredential=recall-api-key:/etc/brain/recall-api-key
ExecStart=/usr/bin/bash /opt/brain/scripts/recall-scheduled-apply.sh
ExecStartPost=/usr/bin/node /opt/brain/scripts/recall-sync-lifecycle-prod.mjs schedule --timer-name brain-recall-sync.timer --db-path /opt/brain/data/brain.sqlite --migrations-dir /opt/brain/scripts/db/migrations
TimeoutStartSec=7h
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=brain-recall-sync

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
RuntimeDirectory=brain-recall-backup-staging
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
ReadWritePaths=/opt/brain/data
ReadWritePaths=/run/brain-recall
ReadWritePaths=/run/brain-recall-backup-staging
UMask=0007
